# rootfs/etc/cron.d/whmcs
*/5 * * * * php -q /var/www/client/crons/cron.php
*/5 * * * * php -q /var/www/client/crons/pop.php

