# rootfs/etc/cron.d/whmcs
0 0 * * * root /usr/bin/php -q /var/www/html/client/crons/cron.php

