# rootfs/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon = true

[program:php7.0-fpm]
command=service php7.0-fpm start

[program:nginx]
command=nginx

[program:cron]
command=cron -f
